# supervisord.conf
[supervisord]
nodaemon=true
logfile=/tmp/supervisord.log

[program:php-fpm7]
command=php-fpm7 -F
stdout_logfile=/dev/stdout
stdout_logfile_maxbytes = 0
stderr_logfile=/dev/stderr
stderr_logfile_maxbytes=0

[program:lighttpd]
command=lighttpd -D -f /etc/lighttpd/lighttpd.custom.conf
stdout_logfile=/dev/stdout
stdout_logfile_maxbytes = 0
stderr_logfile=/dev/stderr
stderr_logfile_maxbytes=0
